# config/upstream/vertex.conf
description "node app {{appid}}+{{hostname}}"

start on started networking
stop on runlevel [016]

script
        cd '/opt/vertex'
        npm start
end script
